# config/control.cron
* * * * * root /usr/local/uravo/bin/control.pl&

