# config/outpost.cron
* * * * * root /usr/local/uravo/bin/outpost.pl&

